# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

; 1. Virtual framebuffer — must start first
[program:xvfb]
command=Xvfb :1 -screen 0 %(ENV_VNC_RESOLUTION)s
autorestart=true
priority=10
stdout_logfile=/var/log/xvfb.log
stderr_logfile=/var/log/xvfb.log

; 2. LXDE desktop — waits briefly for Xvfb
[program:lxde]
command=/bin/bash -c "sleep 2 && lxsession -s LXDE -e LXDE"
environment=DISPLAY=:1,HOME=/root,USER=root
autorestart=true
priority=20
stdout_logfile=/var/log/lxde.log
stderr_logfile=/var/log/lxde.log

; 3. x11vnc — waits for LXDE to be ready
[program:x11vnc]
command=/bin/bash -c "sleep 4 && x11vnc -display :1 -rfbauth /root/.vnc/passwd -rfbport 5900 -forever -shared -noxdamage -wait 50"
autorestart=true
priority=30
stdout_logfile=/var/log/x11vnc.log
stderr_logfile=/var/log/x11vnc.log

; 4. noVNC websocket proxy — browser access on port 6080
[program:novnc]
command=/bin/bash -c "sleep 5 && websockify --web /usr/share/novnc/ 6080 localhost:5900"
autorestart=true
priority=40
stdout_logfile=/var/log/novnc.log
stderr_logfile=/var/log/novnc.log
